# legacy/com.dothax.hl-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dothax.hl-tracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shroperu/Documents/workspace/trader_position_analytics/.venv/bin/python3.14</string>
        <string>/Users/shroperu/Documents/workspace/trader_position_analytics/scripts/run_all.py</string>
        <string>--excel</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shroperu/Documents/workspace/trader_position_analytics</string>
    <key>StandardOutPath</key>
    <string>/Users/shroperu/Documents/workspace/trader_position_analytics/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shroperu/Documents/workspace/trader_position_analytics/logs/launchd_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>VIRTUAL_ENV</key>
        <string>/Users/shroperu/Documents/workspace/trader_position_analytics/.venv</string>
    </dict>
</dict>
</plist>
